# terraform/fabio/config/fabio.conf
description "Fabio"
start on started networking
stop on runlevel [!2345]

respawn
# This is to avoid Upstart re-spawning the process upon `consul leave`
normal exit 0 INT

exec /usr/local/bin/fabio
